# systemd/tusabot-web.service
[Unit]
Description=TusaBot Web Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tusabot/project
ExecStart=/usr/bin/npx serve -s dist -l 5173
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target